# states/supervisor/generic.conf.jinja
[program:{{ name }}]
directory={{ cwd }}
user={{ user }}

command={{ command }}
stdout_logfile={{ log_path }}
redirect_stderr=true
stopsignal=INT
autorestart=true
